# private_dot_config/systemd/user/streamdeck.path
[Unit]
Description="Stream Deck Device Path"

[Path]
PathExists=/dev/streamdeck
Unit=streamdeck.service

[Install]
WantedBy=default.target
